# systemd/ant-node.service
[Unit]
Description=Ant Node - Quantum-proof decentralized network node
Documentation=https://github.com/WithAutonomi/ant-node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ant
Group=ant
ExecStart=/usr/bin/ant-node --config /etc/ant/config.toml --stop-on-upgrade
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s
TimeoutStopSec=30s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/lib/ant /var/log/ant

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
